# scripts/dauerwelle-server.service
[Unit]
Description=Dauerwelle Screen Server
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=david
WorkingDirectory=/home/david/server
ExecStart=/usr/bin/tmux -S /tmp/dauerwelle.sock new-session -d -s dauerwelle 'node /home/david/server/server.js'
ExecStop=/usr/bin/tmux -S /tmp/dauerwelle.sock kill-session -t dauerwelle
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
